# systemd/xdpguard.service
[Unit]
Description=XDPGuard DDoS Protection System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/xdpguard
ExecStart=/usr/bin/python3 /opt/xdpguard/daemon.py
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=false
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
